# systemd/facelock-daemon.service
[Unit]
Description=Facelock Face Authentication Daemon
After=local-fs.target

[Service]
Type=dbus
BusName=org.facelock.Daemon
ExecStart=/usr/bin/facelock daemon
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=3
LimitNOFILE=1024

[Install]
WantedBy=multi-user.target
