# systemd/howdy-daemon.service
[Unit]
Description=Howdy Face Authentication Daemon
Documentation=https://github.com/user/howdy-rust
After=local-fs.target

[Service]
Type=notify
ExecStart=/usr/bin/howdy-daemon
Restart=on-failure
RestartSec=3

# Runtime directory (created automatically with correct ownership)
RuntimeDirectory=howdy
RuntimeDirectoryMode=0755

# Filesystem isolation
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/howdy /run/howdy /var/log/howdy
PrivateTmp=yes

# Device access (cameras only)
DeviceAllow=/dev/video* rw
DevicePolicy=closed

# Security hardening
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
SystemCallFilter=@system-service @io-event
SystemCallArchitectures=native
ProtectProc=invisible
ProcSubset=pid

# Resource limits
LimitNOFILE=1024
LimitNPROC=64

[Install]
WantedBy=multi-user.target
