# systemd/visage-daemon.service
[Unit]
Description=Visage Face Authentication Daemon
Documentation=https://github.com/user/visage
After=local-fs.target visage-daemon.socket
Requires=visage-daemon.socket

[Service]
Type=simple
ExecStart=/usr/bin/visage daemon
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=3

# Resource limits
LimitNOFILE=1024

[Install]
WantedBy=multi-user.target
