# dist/runit/log/run
#!/bin/sh
# runit log script for facelock daemon
exec svlogd -tt /var/log/facelock
